# data/doorpi.service.in
[Unit]
Description=!!project!!
Documentation=https://github.com/motom001/DoorPi/wiki
Requires=doorpi.socket
After=doorpi.socket

[Service]
Type=notify
ExecStart=!!prefix!!/bin/doorpi --configfile !!cfgfile!!
Restart=always
NotifyAccess=all
WatchdogSec=15s

User=!!user!!
Group=!!group!!
SupplementaryGroups=audio
SupplementaryGroups=gpio
SupplementaryGroups=spi
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target
